# setup/com.agent-hud.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-hud.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/haquefarazul/Projects/Agent-Manager/packages/daemon/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/haquefarazul/Projects/Agent-Manager/packages/daemon</string>

    <!-- Start on login and keep it alive (restart if it ever exits). -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Absolute tool paths are used in code, but give launchd a sane PATH too
         so node, caffeinate, pmset and sudo all resolve. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/haquefarazul/Library/Logs/agent-hud/daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/haquefarazul/Library/Logs/agent-hud/daemon.err.log</string>

    <!-- Throttle restart loops if it crashes immediately. -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
